# childrens-stories-worker.service
[Unit]
Description=Children's Stories ARQ Worker
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=kavin
WorkingDirectory=/home/kavin/childrens_stories
Environment="PATH=/home/kavin/.cache/pypoetry/virtualenvs/childrens-stories-6rIxPkgJ-py3.12/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/kavin/.local/bin/poetry run arq backend.worker.WorkerSettings
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
